# etc/systemd/system/update-unbound-ipv6.service
[Unit]
Description=Detects IPv6 prefix changes and rewrites local Unbound config.
Documentation=man:unbound.conf(5)
After=network-online.target unbound.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0027

ExecStartPre=/bin/sh -ec 'test -x /usr/local/bin/update-unbound-ipv6.sh'

Environment=BACKUP_DIR=/var/backups/update-unbound-ipv6
Environment=STATUS_DIR=/var/lib/update-unbound-ipv6
Environment=WEBROOT_DIR=/var/www/html
Environment=LOCK_DIR=/var/lock/update-unbound-ipv6.lock

ExecStart=/usr/local/bin/update-unbound-ipv6.sh

StandardOutput=journal
StandardError=journal

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

ReadWritePaths=/etc/unbound/unbound.conf.d /var/backups/update-unbound-ipv6 /var/lib/update-unbound-ipv6 /var/www/html /var/lock

[Install]
WantedBy=multi-user.target
