# etc/systemd/system/update-unbound-ipv6.service
[Unit]
Description=Dynamic IPv6 config rewriting, validation, backups.
Documentation=man:unbound.conf(5)
After=network-online.target unbound.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0027

ExecStartPre=/usr/bin/install -d -m 0750 -o root -g root /var/backups/update-unbound-ipv6

ExecStartPre=/bin/sh -ec 'for b in awk basename cat chmod chown cp cut date find grep head ip mkdir mktemp mv rm sed sort stat systemctl tail tr unbound unbound-checkconf xargs; do command -v "$b" >/dev/null 2>&1 || { printf "Missing required binary: %s\n" "$b" >&2; exit 1; }; done'
ExecStartPre=/bin/sh -ec 'test -x /usr/local/bin/update-unbound-ipv6.sh'

Environment=BACKUP_DIR=/var/backups/update-unbound-ipv6

ExecStart=/usr/local/bin/update-unbound-ipv6.sh

StandardOutput=journal
StandardError=journal

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

ReadWritePaths=/etc/unbound /var/backups

[Install]
WantedBy=multi-user.target
