# etc/systemd/system/update-unbound-ipv6.service
[Unit]
Description=Detects IPv6 prefix changes and rewrites local Unbound config.
Documentation=man:unbound.conf(5)
After=network-online.target unbound.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0027

ExecStartPre=/usr/bin/install -d -m 0750 /var/backups/update-unbound-ipv6
ExecStartPre=/bin/sh -ec 'test -x /usr/local/bin/update-unbound-ipv6.sh'
ExecStart=/usr/local/bin/update-unbound-ipv6.sh

StateDirectory=update-unbound-ipv6
RuntimeDirectory=update-unbound-ipv6
RuntimeDirectoryMode=0750

Environment=STATUS_DIR=%S/update-unbound-ipv6
Environment=BACKUP_DIR=/var/backups/update-unbound-ipv6
Environment=WEBROOT_DIR=/var/www/html
Environment=LOCK_DIR=%t/update-unbound-ipv6/lock

StandardOutput=journal
StandardError=journal

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

ReadWritePaths=/etc/unbound/unbound.conf.d /var/backups /var/www/html

[Install]
WantedBy=multi-user.target
